# tools/frr@.service.in
[Unit]
Description=FRRouting
Documentation=https://frrouting.readthedocs.io/en/latest/setup.html
Wants=network.target
After=network-pre.target systemd-sysctl.service
Before=network.target
OnFailure=heartbeat-failed@%n

[Service]
Nice=-5
Type=forking
NotifyAccess=all
StartLimitInterval=3m
StartLimitBurst=3
TimeoutSec=@TIMEOUT_MIN@m
WatchdogSec=60s
RestartSec=5
Restart=always
LimitNOFILE=1024
<<<<<<< HEAD
PIDFile=@CFG_STATE@/%I/watchfrr.pid
ExecStart=@CFG_SBIN@/frrinit.sh start %I
ExecStop=@CFG_SBIN@/frrinit.sh stop %I
ExecReload=@CFG_SBIN@/frrinit.sh reload %I
=======
PIDFile=@e_frr_runstatedir@/%I/watchfrr.pid
ExecStart=@e_sbindir@/frrinit.sh start %I
ExecStop=@e_sbindir@/frrinit.sh stop %I
ExecReload=@e_sbindir@/frrinit.sh reload %I
>>>>>>> 9b0b9282d (bgpd: Fix bgp core with a possible Intf delete)

[Install]
WantedBy=multi-user.target
